# files/etc/cron.d/cleanup-youtube
0 * * * * /usr/local/bin/cleanup-youtube.sh /logs/cron.log 2>&1


