# files/etc/cron.hourly/check-twitter
#!/bin/bash

{
    /usr/local/bin/twitter-bot -c /secret/discord.conf check-feed
} >> /logs/twitter-bot.log 2>&1


