# files/etc/cron.hourly/logrotate
#!/bin/bash

/usr/sbin/logrotate -s /var/run/logrotate.pid /etc/logrotate.conf
